# pisugar-weather.service
[Unit]
Description=PiSugar Weather Display
After=network.target

[Service]
Type=simple
User=waboring
WorkingDirectory=/home/waboring/pisugar-wx
ExecStart=/usr/bin/python3 -m src.app
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
