# pisugar-weather.service
[Unit]
Description=PiSugar Weather Display
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/pisugar-wx
ExecStart=/home/pi/pisugar-wx/venv/bin/python -m src.app
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
